# deploy/install/templates/baluhost-backend-local.socket
[Unit]
Description=BaluHost local-channel socket (Tauri Companion)

[Socket]
ListenStream=/run/baluhost/local.sock
SocketMode=0660
SocketUser=baluhost
SocketGroup=baluhost
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
