# deploy/install/templates/baluhost-backend.service
[Unit]
Description=BaluHost NAS Backend (FastAPI/Uvicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=@@BALUHOST_USER@@
Group=baluhost
WorkingDirectory=@@INSTALL_DIR@@/backend

# Load production environment variables
EnvironmentFile=@@INSTALL_DIR@@/.env.production

# Clean up stale primary worker lock before start
ExecStartPre=/usr/bin/rm -f /tmp/baluhost-primary.lock

# Virtual environment activation and uvicorn start.
# --proxy-headers + --forwarded-allow-ips=127.0.0.1: trust X-Forwarded-For ONLY
# from the local nginx upstream, so request.client.host is the real client IP
# (per-client rate limiting + accurate audit IPs). The trust pin to 127.0.0.1
# prevents X-Forwarded-For spoofing, incl. direct connections to :8000.
ExecStart=@@VENV_BIN@@/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --log-level info \
    --no-access-log \
    --proxy-headers \
    --forwarded-allow-ips=127.0.0.1

# Restart policy
Restart=always
RestartSec=10s

# Security hardening
# NoNewPrivileges must remain false — backend needs sudo for mdadm, hwmon/PWM, cpufreq
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=baluhost-backend

[Install]
WantedBy=multi-user.target
