# deploy/systemd/baluhost-backend.service
[Unit]
Description=BaluHost NAS Backend (FastAPI/Uvicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=sven
Group=sven
WorkingDirectory=/home/sven/projects/BaluHost/backend

# Load production environment variables
EnvironmentFile=/home/sven/projects/BaluHost/.env.production

# Virtual environment activation and uvicorn start
# 4 workers: primary detection via file lock (/tmp/baluhost-primary.lock).
# Hardware services (fan, power, telemetry) run only on the primary worker.
# Nginx ip_hash ensures per-client sticky sessions across workers.
ExecStartPre=/usr/bin/rm -f /tmp/baluhost-primary.lock
ExecStart=/home/sven/projects/BaluHost/backend/.venv/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --log-level info \
    --no-access-log

# Fallback: notify BaluPi before SIGTERM (idempotent, harmless if Pi is offline or disabled).
# The primary notification happens inside the FastAPI lifespan shutdown handler.
# This curl fires before SIGTERM so it's race-free even if the app doesn't shut down cleanly.
# Timeout 5s to avoid blocking systemd. Failure is silently ignored (-sf).
ExecStop=/bin/sh -c 'curl -sf -m 5 -X POST http://${BALUPI_URL:-localhost:9999}/api/handshake/nas-going-offline -H "Content-Type: application/json" -d "{}" 2>/dev/null || true'

# Restart policy
Restart=always
RestartSec=10s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=baluhost-backend

[Install]
WantedBy=multi-user.target
